# deploy/systemd/bin-core.service
[Unit]
Description=Bin Trading Core
After=network.target postgresql.service redis.service

[Service]
WorkingDirectory=/home/ruthgyeul/TossInvestAI
ExecStart=/home/ruthgyeul/TossInvestAI/venv/bin/python3 -m core.main
Restart=always
RestartSec=10
Environment=TZ=Asia/Seoul

[Install]
WantedBy=multi-user.target
